# deploy/spotonauto-web.service
[Unit]
Description=SpotOnAuto Next.js Web App
After=network-online.target postgresql@16-main.service lemon-manuals.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/spotonauto.com
EnvironmentFile=/root/spotonauto.com/.env.local
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=127.0.0.1
ExecStartPre=/usr/bin/npm run prestart
ExecStart=/usr/bin/npm run start:web
Restart=always
RestartSec=5
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
